# eBite/supervisor/LeverHandleQueue.conf
#处理单子
[program:queue-userLeverHandle]
process_name=%(process_num)02d
directory=/www/wwwroot/m.mddz.cc
command=php artisan queue:work --queue=lever:handle --sleep=1 --tries=3
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/null
