# eBite/supervisor/microUpdateQueue.conf
#秒合约平仓队列
[program:queue-microUpdate]
process_name=%(process_num)02d
directory=/www/wwwroot/m.mddz.cc
command=php artisan queue:work --queue=micro:update:%(process_num)02d --sleep=1 --tries=3
autostart=true
autorestart=true
user=www
#只能开一个进程,否则会发生重复平仓
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/null
