# eBite/supervisor/pushMarket.conf
#推送行情队列
[program:queue-pushMarket]
process_name=%(program_name)s_%(process_num)02d
directory=/www/wwwroot/m.mddz.cc
command=php artisan queue:work --queue=pushMarket  --sleep=1 --tries=3
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/null
